# roles/postfix/templates/postfix.conf.j2
[program:postfix]
command=/bin/bash -c "postfix start"
user=postfix
autostart=true
autorestart=unexpected
