# roles/postfix/templates/postfix.conf.j2
[program:postfix]
command=/usr/sbin/postfix
user=postfix
autostart=true
autorestart=unexpected
