# roles/girder/templates/girderd.conf.j2
[program:girder]
command=/opt/girder/env/bin/python -m girder
directory=/opt/girder
user={{ grits_user }}
autostart=true
autorestart=unexpected
environment=GIRDER_MONGO_URI="{{ mongo_url }}/girder"
